# haosap/haosap/rootfs/etc/services.d/haosap/run
#!/bin/bash

exec /haosap.sh  >> /var/log/haosap.log 2>&1
