# cloudflare_ip_updater/rootfs/etc/services.d/cloudflare/finish
#!/bin/bash

echo "Cloudflare service stopping" >> /var/log/cloudflare.log

# Add any additional cleanup commands here

exit 0
