# cloudflare_ip_updater/rootfs/etc/services.d/cloudflare/run
#!/bin/bash

exec /cloudflare.sh  >> /var/log/cloudflare.log 2>&1
